# components/congress/ansible/roles/deploy/templates/congress-api.conf
description "Congress API server"

env PYTHON_PATH={{installDir}}/congress
start on runlevel [2345]
stop on runlevel [!2345]

exec {{installDir}}/congress/bin/congress-server --config-file /etc/congress/congress.conf
